# system/alienware-fan/alienware-fan.service
[Unit]
Description=Alienware m16 R2 fan control
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/alienware-fan
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
